# examples/dbgpproxy_supervisord.conf
[program:dbgpproxy]
command=/usr/src/dbgpproxy/bin/dbgpproxy -i 192.168.33.114:9001 -d 127.0.0.1:9000
exitcodes=0
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/tmp/dbgpproxy.log